# deploy/systemd/careos-lite-scheduler.service
[Unit]
Description=CareOS Lite Scheduler Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/careos-lite
EnvironmentFile=/opt/careos-lite/.env
ExecStart=/opt/careos-lite/.venv/bin/python -m careos.workers.scheduler_worker
Restart=always
RestartSec=3
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
