# deploy/systemd/careos-lite-scheduler.service
[Unit]
Description=CareOS Lite Scheduler Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/kumarmankala/careos
EnvironmentFile=/home/kumarmankala/careos/.env
ExecStart=/home/kumarmankala/careos/.venv/bin/python -m careos.workers.scheduler_worker
Restart=always
RestartSec=3
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
